# deploy/agent-life-space.service
[Unit]
Description=Agent Life Space
Documentation=https://github.com/B2JK-Industry/Agent_Life_Space
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# CHANGE these to match your install.
User=youruser
Group=youruser
WorkingDirectory=/home/youruser/Agent_Life_Space
EnvironmentFile=/home/youruser/Agent_Life_Space/.env

# Runtime data + log dirs: configured in .env (EnvironmentFile above).
# Do NOT set AGENT_DATA_DIR / AGENT_LOG_DIR here — .env is the single
# source of truth so that agent reporting matches the actual paths.

# Enable Telegram self-update + automatic post-update restart.
# When the operator messages "nasad novú verziu u seba" the agent
# will pull the latest fast-forward, drain in-flight work, exit 0,
# and this unit will bring up a fresh process with the new code.
Environment=AGENT_SELF_RESTART_AFTER_UPDATE=1

# Use the project venv.
ExecStart=/home/youruser/Agent_Life_Space/.venv/bin/python -m agent

# Restart policy: bring us back on every clean exit (including the
# self-restart after a successful update). RestartSec keeps the gap
# short; StartLimitBurst stops infinite loops on a bad update.
Restart=always
RestartSec=2
StartLimitBurst=5
StartLimitIntervalSec=60

# Watchdog disabled — the agent runs long build jobs (codegen + Docker)
# that can take 5+ minutes. Re-enable once sd_notify heartbeat is
# implemented in the main loop.
# WatchdogSec=120

# Hardening — adjust to your environment.
# NoNewPrivileges disabled because Docker requires setgid for group
# switching when running containers from the build pipeline.
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
